# systemd/certbot.service.d/cert_renew.conf
# Config to work with 'cert_renew.target'.
#[Service]
#ExecStart=/bin/sh -c 'echo /usr/bin/certbot -q renew "$(date)" >> /root/1.tmp'

[Install]
WantedBy=cert_renew.target
